# crypt.service
[Service]
ExecStart = /home/sanveer/crypto_service/crypt.sh
