# decrypt.service
[Service]
ExecStart = /home/sanveer/crypto_service/decrypt.sh
